# x13s-touchscreen-udev/72-x13s-touchscreen.rules
ACTION=="add" \
, RUN+="/bin/bash -c 'echo 1-0010 > /sys/bus/i2c/drivers/i2c_hid_of/bind'"
